# prod/labour-bureau-compose.service
# Install: sudo cp labour-bureau-compose.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now labour-bureau-compose.service
#
# Edit these before installing (if your setup differs):
#   - WorkingDirectory + EnvironmentFile: path to prod dir (e.g. /home/YOUR_USER/projects/labour-bureau/prod)
#   - User + Group: user that runs podman (run 'podman ps' without sudo; that user)
#   - XDG_RUNTIME_DIR: /run/user/$(id -u) for that user (e.g. /run/user/1000)

[Unit]
Description=Labour Bureau stack via podman compose (Politburo, Comrade Bot, DB, Redis, Monitoring)
Documentation=https://github.com/infinite-experiment
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=simple
# Run compose in foreground so systemd tracks the stack; Restart=always brings it back if it exits
ExecStart=/usr/bin/podman compose -f docker-compose.prod.yml up
ExecStop=/usr/bin/podman compose -f docker-compose.prod.yml down -t 30
TimeoutStopSec=90

WorkingDirectory=/home/runner/projects/labour-bureau/prod
EnvironmentFile=-/home/runner/projects/labour-bureau/prod/env/monitoring.env

User=runner
Group=runner
Environment=XDG_RUNTIME_DIR=/run/user/1000

Restart=always
RestartSec=10
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
