# prod/labour-bureau-services.service
[Unit]
Description=Labour Bureau Services (Politburo, Comrade Bot, DB, Redis, Monitoring)
Documentation=https://github.com/infinite-experiment
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Change User/Group to match the user that runs podman
# For rootless podman, use your username (e.g., User=odin)
# For rootful podman, use root
# To check: run 'podman ps' - if it works without sudo, use your username
User=odin
Group=odin
WorkingDirectory=/home/odin/projects/infinite-experiment/labour-bureau/prod
ExecStart=/home/odin/projects/infinite-experiment/labour-bureau/prod/start-services.sh
ExecStop=/home/odin/projects/infinite-experiment/labour-bureau/prod/stop-services.sh
StandardOutput=journal
StandardError=journal
# Restart the systemd service if it fails (though containers auto-restart)
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
