# prod/podman-log-shipper.service
[Unit]
Description=Podman Log Shipper - Tails container logs to files for Promtail
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=runner
Group=runner
WorkingDirectory=/home/runner/projects/labour-bureau/prod
ExecStart=/usr/bin/bash /home/runner/projects/labour-bureau/prod/scripts/podman-log-shipper.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/runner"

[Install]
WantedBy=multi-user.target
